# tailscale-firstboot.service
[Unit]
Description=Tailscale First‑Boot Authentication
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tailscale up --ssh
StandardInput=tty
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target